# deploy/launchd-orchestrator/templates/daemon/com.stocksystem.main-service.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stocksystem.main-service</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PROJECT_ROOT__/scripts/start_all.sh</string>
    <string>--foreground</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>UserName</key>
  <string>__RUN_AS_USER__</string>
  <key>GroupName</key>
  <string>__RUN_AS_GROUP__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>__PROJECT_ROOT__/logs/launchd/main-service.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_ROOT__/logs/launchd/main-service.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
    <key>ORCHESTRATOR_LAUNCHD_ENV</key>
    <string>__PROJECT_ROOT__/deploy/launchd-orchestrator/launchd.env</string>
  </dict>
</dict>
</plist>